# cron/cron.d/0hourly
# Run the hourly jobs
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
# MAILTO=""
*/1 * * * * operator run-parts /etc/cron.minutely
01 * * * * operator run-parts /etc/cron.hourly
